# spirometry_backend/test_python.timer
[Unit]
Description=Timer that periodically triggers test_python.service

[Timer]
OnCalendar=*-*-* 11:00:00
Persistent=true
