# ecp-hs/dev/dashboard_scripts/daily_check.timer
[Unit]
Description=Timer that periodically triggers daily_check.service

[Timer]
OnCalendar=*-*-* 11:00:00
Persistent=true